# dot_config/powershell/profile.d/00-aliases.ps1
##
# Useful aliases

# chezmoi aliases
if (Get-Command "chezmoi" -ErrorAction SilentlyContinue) {
    Set-Alias -Name cm -Value chezmoi
}
